# pc/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

; --- Syncthing: p2p connection to the server. GUI bound to localhost only;
; on Windows/Docker Desktop that's still reachable at localhost:8384 from
; the Windows browser for first-time pairing.
[program:syncthing]
command=/usr/bin/syncthing serve --no-browser --no-restart --gui-address=127.0.0.1:8384 --home=/data/syncthing-config
autostart=true
autorestart=true
stopsignal=INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- Moves fully-synced submission folders into the permanent archive.
[program:file_mover]
command=python3 /scripts/file_mover.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
