# server/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

; --- Syncthing: p2p sync with the office PC. GUI bound to localhost only,
; never exposed outside the container/host (see docker-compose.yml).
; autostart/autorestart come from ENABLE_SYNCTHING (normalized to
; SYNCTHING_AUTOSTART by entrypoint.sh) - set ENABLE_SYNCTHING=false in
; .env and this program never launches at all.
[program:syncthing]
command=/usr/bin/syncthing serve --no-browser --no-restart --gui-address=127.0.0.1:8384 --home=/data/syncthing-config
autostart=%(ENV_SYNCTHING_AUTOSTART)s
autorestart=%(ENV_SYNCTHING_AUTOSTART)s
stopsignal=INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- The Flask upload form, served by gunicorn on 127.0.0.1:8000.
; A host-level Nginx (outside this container) reverse-proxies HTTPS to it.
[program:webapp]
command=gunicorn --bind 0.0.0.0:8000 --workers 2 --timeout 120 --access-logfile - --error-logfile - --capture-output app:app
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- Deletes submission folders older than MAX_AGE_DAYS on a loop.
[program:cleanup]
command=/scripts/cleanup_loop.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0